# deploy/systemd/moon-system.service
[Unit]
Description=Moon System Watcher
After=network.target

[Service]
Type=simple
WorkingDirectory={{MOON_SYSTEM_DIR}}
ExecStart=/usr/bin/env cargo run --manifest-path {{MOON_SYSTEM_DIR}}/Cargo.toml -- moon-watch --daemon
Restart=always
RestartSec=5
Environment=RUST_LOG=info

[Install]
WantedBy=multi-user.target
